# src/configure/daemons/gadgetini/data_crawler.service
[Unit]
Description=Gadgetini sensor collector (auto-detects PCB Modbus or legacy ADS1256).
After=redis.service
[Service]
Type=simple
ExecStart=python3 /home/gadgetini/gadgetini/src/exporter/data_crawler.py
Restart=always
RestartSec=5
User=gadgetini
Group=gadgetini
SupplementaryGroups=dialout
[Install]
WantedBy=multi-user.target
