# src/gui/deamon/gadgetini_ui.service
[Unit]
Description=Run npm start for gadgetini web ui service:3001
# pcb_bootstrap must finish first so the control_board.service / data_crawler.service
# branch is decided; only then can the UI report the PCB connection state correctly.
After=network.target pcb_bootstrap.service redis.service

[Service]
Type=simple
WorkingDirectory=/home/gadgetini/gadgetini/src/gui/gadgetini-web
ExecStart=/usr/bin/npm start -- --port 3001
Restart=always

[Install]
WantedBy=multi-user.target