# src/gui/deamon/gadgetini_ui.service
[Unit]
Description=Run npm start for gadgetini web ui service:3001
# pcb_bootstrap이 먼저 끝나야 control_board.service / data_crawler.service 분기가
# 결정된 상태에서 UI가 PCB 연결 상태를 정확히 표시할 수 있다.
After=network.target pcb_bootstrap.service redis.service

[Service]
Type=simple
WorkingDirectory=/home/gadgetini/gadgetini/src/gui/gadgetini-web
ExecStart=/usr/bin/npm start -- --port 3001
Restart=always

[Install]
WantedBy=multi-user.target