# src/configure/daemons/host/gadgetini_net_up.service
[Unit]
Description=Activate gadgetini monitoring system network start up
Requires=NetworkManager.service
After=network-online.target

[Service]
# Bring up the profile created by usb-gadget-host.sh (name was mismatched as
# "gadgetini" before, which never existed → unit failed while autoconnect still
# worked). The profile is bound to its ifname, and `|| true` tolerates it already
# being active via connection.autoconnect.
ExecStart=/bin/sh -c '/usr/bin/nmcli connection up usb-gadget-host || true'
Type=oneshot
RemainAfterExit=yes
User=root
Group=root
[Install]
WantedBy=multi-user.target

